# mbmex.service
[Unit]
Description=Many Button Mouse Extend

[Service]
ExecStart=/usr/local/bin/mbmex
Restart=always

[Install]
WantedBy=multi-user.target
