# mbmex.service
[Unit]
Description=Many Button Mouse Extend

[install]
WantedBy=multi-user.target

[Service]
ExecStart=/usr/local/bin/mbmex
Restart=always
